# Library/LaunchAgents/com.patrick.system-state-check.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.patrick.system-state-check</string>

	<key>ProgramArguments</key>
	<array>
		<string>{{ .chezmoi.homeDir }}/code/scripts-private/bin/system-state-check</string>
		<string>--json</string>
		<string>--no-fetch</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>StartInterval</key>
	<integer>300</integer>

	<key>StandardOutPath</key>
	<string>/dev/null</string>

	<key>StandardErrorPath</key>
	<string>/tmp/system-state-check.error.log</string>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
