# Library/LaunchAgents/com.aerospace.display-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.aerospace.display-monitor</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/pmotard/code/scripts/aerospace-display-monitor</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/tmp/aerospace-display-monitor.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/aerospace-display-monitor.error.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
